# dist/rsync-proxy.service
[Unit]
Description=Rsync Proxy

[Service]
Restart=on-failure
StartLimitInterval=10s

ExecStart=/usr/local/bin/rsync-proxy -v=3 --config=/etc/rsync-proxy/config.toml
ExecReload=/usr/local/bin/rsync-proxy --config=/etc/rsync-proxy/config.toml --reload

User=nobody
Group=nogroup
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
